# ec2mc/aws_setup_src/user_data/mc_template/crontabs/ec2mc_crontab
SHELL=/bin/bash
PATH=/usr/bin/:/usr/sbin/
@reboot root /home/ssm-user/manage-scripts/startup_script.sh > /dev/null
* * * * * root /home/ssm-user/manage-scripts/connection_listener.sh > /dev/null
