# ec2mc/aws_setup_src/user_data/mc_template/crontabs/mc_crontab
@reboot root ~/manage-scripts/startup_script.sh
* * * * * root ~/manage-scripts/connection_listener.sh